# etc/systemd/system/makeswap-on-azure.service
# Required by systemd for interpretation of service operation.

[Unit]
Description=Makes a swap file on the volitile temp drive

[Service]
<<<<<<< HEAD
<<<<<<< Updated upstream
ExecStart=/bin/bash -c 'chmod ug+rx /usr/local/bin/makeswap.sh'
ExecStart=/usr/local/bin/makeswap.sh
=======
ExecStartPre=+chmod ug+rx /usr/sbin/makeswap.sh
ExecStart=+/usr/sbin/makeswap.sh
>>>>>>> Stashed changes
=======
ExecStartPre=+chmod ug+rx /usr/local/bin/makeswap.sh
ExecStart=+/usr/local/bin/makeswap.sh
>>>>>>> 6bc2c1c32770c0c4b6f556048c9e4f302c7241bd
Type=oneshot
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target