# etc/systemd/system/makeswap-on-azure.service
# Required by systemd for interpretation of service operation.

[Unit]
Description=Makes a swap file on the volitile temporary drive

[Service]
ExecStartPre=+chmod ug+rx /usr/local/bin/makeswap.sh
ExecStart=+/usr/local/bin/makeswap.sh
Type=oneshot
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target